# cron/pet-energy.timer
[Unit]
Description=Run pet energy decay every 20 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=20min
Persistent=true

[Install]
WantedBy=timers.target
